# systemd/organia-frontend-demo.service
[Unit]
ExecStartPre=podman rm -f %N
ExecStart=podman run --rm --name %N -p ${PORT}:8080 organia-frontend:demo